# OPTB_daily_noaa_fetch.timer
[Unit]
Description=Trigger the Ocean Park Tides Bot daily NOAA data fetch daily at 3am.

[Timer]
OnCalendar=03:00

[Install]
WantedBy=timers.target
